# backend/supervisord.conf
[program:exporter]
command = perl /home/dotcloud/current/bin/exporter

[program:saver]
command = perl /home/dotcloud/current/bin/saver
